# crontab
0 0 * * * /root/update.sh >/dev/tty 2>&1
