# crontab
0 0 * * * root /root/update.sh >/dev/tty 2>&1
