# crontab
0 0 * * * /bin/bash /root/update.sh >> /var/log/cron.log 2>&1
